# deploy_tools/gunicorn-systemd.service
[Unit]
Description=Gunicorn server for SITENAME

[Service]
Restart=on-failure
User=aries
WorkingDirectory=/home/aries/sites/SITENAME/source
ExecStart=/home/aries/sites/SITENAME/virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket superlists.wsgi:application

[Install]
WantedBy=multi-user.target